# sample_program.conf
[program:sample]
directory=/apps/www/sample
command=java sample.jar
autostart=true
autorestart=false
startretries=10
redirect_stderr=true
stdout_logfile=/var/log/sample.log
